# image/rootfs/var/lib/opkg/info/busybox.list
/etc/busybox.links
/bin/sh
/etc/init.d/inetd.busybox
/bin/busybox
/etc/inetd.conf
