# image/rootfs/var/lib/opkg/info/busybox-syslog.list
/etc/syslog-startup.conf.busybox
/etc/default/busybox-syslog
/etc/init.d/syslog.busybox
