# image/rootfs/var/lib/opkg/info/init-ifupdown.list
/etc/network/if-pre-up.d/nfsroot
/etc/init.d/networking
/etc/network/interfaces
/etc/network/if-up.d/
