# image/rootfs/var/lib/opkg/info/matrix-gui.list
/usr/share/matrix-gui-2.0/css/
/usr/share/matrix-gui-2.0/css/fonts-min.css
/usr/share/matrix-gui-2.0/css/800x480.css
/usr/share/matrix-gui-2.0/images/exit-icon.png
/usr/share/matrix-gui-2.0/submenu.php
/usr/share/matrix-gui-2.0/css/800x600.css
/usr/share/matrix-gui-2.0/images/right-arrow-icon.png
/usr/share/matrix-gui-2.0/images/
/usr/share/matrix-gui-2.0/images/tex.png
/usr/share/matrix-gui-2.0/app_description.php
/etc/init.d/matrix-gui-2.0
/usr/share/matrix-gui-2.0/coming_soon.php
/usr/share/matrix-gui-2.0/images/left-arrow-icon.png
/usr/share/matrix-gui-2.0/images/run-icon.png
/usr/share/matrix-gui-2.0/css/yui-reset-fonts-base.css
/usr/share/matrix-gui-2.0/images/multi-icon.png
/usr/share/matrix-gui-2.0/javascript/jquery-latest.js
/usr/share/matrix-gui-2.0/Makefile
/usr/share/matrix-gui-2.0/css/global.css
/usr/share/matrix-gui-2.0/javascript/
/usr/share/matrix-gui-2.0/php.ini
/usr/share/matrix-gui-2.0/menubar.php
/usr/share/matrix-gui-2.0/css/640x480.css
/usr/share/matrix-gui-2.0/run_script.php
/usr/share/matrix-gui-2.0/execute_command.sh
/usr/share/matrix-gui-2.0/index.php
/usr/share/matrix-gui-2.0/execute-command.php
/usr/share/matrix-gui-2.0/LICENSE
/usr/share/matrix-gui-2.0/css/480x272.css
/usr/share/matrix-gui-2.0/generate.php
/usr/share/matrix-gui-2.0/supported_resolutions.txt
/usr/share/matrix-gui-2.0/images/coming-icon.png
/usr/share/matrix-gui-2.0/images/back-arrow-icon.png
/usr/share/matrix-gui-2.0/images/down-arrow-icon.png
/usr/share/matrix-gui-2.0/helper_functions.php
