# image/rootfs/var/lib/opkg/info/ts-calibrate-init.postinst
#!/bin/sh
if test "x$D" != "x"; then
	OPT="-r $D"
else
	OPT="-s"
fi
if type update-rc.d >/dev/null 2>/dev/null; then
	update-rc.d $OPT ts-calibrate start 96 2 3 4 5 .
fi
